# contrib/init/paladeumd.conf
description "Paladeum Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env PLBD_BIN="/usr/bin/paladeumd"
env PLBD_USER="paladeum"
env PLBD_GROUP="paladeum"
env PLBD_PIDDIR="/var/run/paladeumd"
# upstart can't handle variables constructed with other variables
env PLBD_PIDFILE="/var/run/paladeumd/paladeumd.pid"
env PLBD_CONFIGFILE="/etc/paladeum/paladeum.conf"
env PLBD_DATADIR="/var/lib/paladeumd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # paladeumd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$PLBD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run paladeumd."
        echo "The setting must appear in $PLBD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Paladeum Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$PLBD_PIDDIR"
    chmod 0755 "$PLBD_PIDDIR"
    chown $PLBD_USER:$PLBD_GROUP "$PLBD_PIDDIR"
    chown $PLBD_USER:$PLBD_GROUP "$PLBD_CONFIGFILE"
    chmod 0660 "$PLBD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$PLBD_PIDFILE" \
    --chuid $PLBD_USER:$PLBD_GROUP \
    --exec "$PLBD_BIN" \
    -- \
    -pid="$PLBD_PIDFILE" \
    -conf="$PLBD_CONFIGFILE" \
    -datadir="$PLBD_DATADIR" \
    -disablewallet \
    -daemon

